# server/script/service/backend.service
[Unit]
Description=Node.js Backend Service

[Service]
WorkingDirectory=/home/STRAMATEL/server/backend
ExecStart=/home/STRAMATEL/server/script/start_frontend.sh
User=STRAMATEL
Group=STRAMATEL
Restart=always

[Install]
WantedBy=multi-user.target

